# linux_lab1/lab.service
[Unit]
Description=lab demo target

[Service]
Type=simple
ExecStart=/usr/bin/whatever.sh 
RemainAfterExit=yes

[Install]
WantedBy=graphical.target
